# package/openwrt/net/netatalk/files/afpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2012 OpenWrt.org

START=70


start()
{
	mkdir -p /var/netatalk/msg
	mkdir -p /var/netatalk/CNID
	service_start /usr/sbin/netatalk
}

stop()
{
	service_stop /usr/sbin/netatalk
}

reload()
{
	service_reload /usr/sbin/netatalk
}
